# bin/tunnerse-server.service
[Unit]
Description=Tunnerse Server - Local tunnel management daemon
After=network.target

[Service]
Type=simple
User=pedroborgezs
Group=pedroborgezs
ExecStart=/usr/local/bin/tunnerse-server
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
